# imageFiles/etc/sv/syslog-ng/run
#!/usr/bin/env sh
set -eu
exec 2>&1

/usr/sbin/syslog-ng --foreground
